# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:backend]
command=node server.js
directory=/app/backend
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true